# deploy/systemd/openclaw.service
[Unit]
Description=OpenClaw Agent Runtime
Documentation=https://github.com/Eldonlandsupply/EldonOpenClaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw

# Edit this to your actual install path
WorkingDirectory=/opt/openclaw

# Secrets file — chmod 600, owned by openclaw
# Create with: sudo cp .env.example /etc/openclaw/openclaw.env && sudo chmod 600 /etc/openclaw/openclaw.env
EnvironmentFile=/etc/openclaw/openclaw.env

ExecStart=/opt/openclaw/.venv/bin/python -m openclaw.main /opt/openclaw/config.yaml

# 24/7 restart policy
Restart=always
RestartSec=5
TimeoutStopSec=20

# Structured JSON logs go to journald automatically
# View with: journalctl -u openclaw -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/openclaw

[Install]
WantedBy=multi-user.target
